# source/etc/systemd/system/audio-teleport-receiver.socket
[Unit]
Description=Audio Teleport Receiver Socket

[Socket]
ListenStream=0.0.0.0:12345
Accept=yes

[Install]
WantedBy=sockets.target
